# docket/deploy/owl-docket-agent.service
[Unit]
Description=Docket autonomous agent - orchestrator that works tickets off the queue
After=network-online.target owl-docket.service
Wants=network-online.target

[Service]
Type=simple
# Runs as root: has claude creds (/root/.claude), the GitHub push key
# (/root/.ssh), and the Neil B <thenofisamizdat@gmail.com> commit identity.
User=root
Group=root
WorkingDirectory=/home/conorbowles51/app_v2-docket/backend
Environment="HOME=/root"
Environment="PATH=/root/.local/bin:/home/conorbowles51/app_v2/venv/bin:/usr/local/bin:/usr/bin:/bin"
# Autonomy gate: 0 = grooming only (assess + plan, read-only). Set to 1 to let
# the agent implement, self-review and open PRs.
Environment="DOCKET_AGENT_WRITES=1"
Environment="DOCKET_AGENT_MODEL=sonnet"
# Real PR objects: set a GitHub token (repo scope) and the agent opens PRs via
# the API instead of recording a compare URL. Uncomment when the PAT lands.
#Environment="DOCKET_GITHUB_TOKEN=ghp_..."
# Email notifications: drained via msmtp once /etc/msmtprc has the SMTP cred
# (see msmtprc.example next to this unit). Sender identity:
#Environment="DOCKET_MAIL_FROM=Docket <docket@example.com>"
ExecStart=/home/conorbowles51/app_v2/venv/bin/python -m services.docket_agent
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
